# supervisor/subscription.conf
[program:subscription]
command = /home/pi/go/bin/subscription -sample=false -addr=:80 -config=/home/pi/subscriptions.json
autostart = false
autorestart = unexpected
startretries = 3
